# camera.service
[Unit]
Description=Camera streaming service
After=network.target

[Service]
Type=simple
User=roost
WorkingDirectory=/home/roost/camera
ExecStart=/home/roost/camera/target/release/camera
Restart=on-failure
RestartSec=2

[Install]
WantedBy=multi-user.target
